# etc/qnewsflash.service.in
[Unit]
Description=Server utilities for qNewsFlash
Before=network-online.target
After=NetworkManager.service systemd-resolved.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 30
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/qNewsFlash --sysconf=@SYSCONF_DIR@/qNewsFlash.json
ExecReload=/bin/kill -HUP  $MAINPID
ExecStop=/bin/kill   -TERM $MAINPID
PIDFile=@PID@/qNewsFlash.pid
Restart=on-success

[Install]
WantedBy=multi-user.target
