# etc/qnewsflash.timer.in
[Unit]
Description=Timer for qNewsFlash
Before=network-online.target
After=NetworkManager.service systemd-resolved.service

[Timer]
# システム起動後30秒後にqnewsflash.serviceサービスを起動する
OnBootSec=30
Unit=qnewsflash.service

[Install]
# システムのタイマターゲットに関連付け、システム起動時にタイマを有効にする
WantedBy=timers.target
